# crontab
# Put the following lines in the crontab file after using the command: crontab -e
@reboot ( sleep 15; sh /usr/local/bin/pseyestartscript & )